# supervisor/supervisor.d/cgihttpserver.ini
[program:cgihttpserver]
command = /usr/bin/python  -m  CGIHTTPServer   10000 
directory = /home/work 
